# haproxy.cron
#!/bin/bash

#
# Renew certificates from LetsEncrypt
#

if [ -f /.env-haproxy-letsencrpyt ]; then
	/usr/local/sbin/certbot-renew >/dev/null
fi


#
# Reload any configuration changes
#

if [ -f /.env-haproxy-daily-reload ]; then
        systemctl reload haproxy
fi


#
# Wait for reload to finish then update OCSP stapling
#

if [ -f /.env-haproxy-letsencrpyt ] && [ -f /.env-haproxy-letsencrpyt-ocsp ]; then
	sleep 15
	/usr/local/sbin/certbot-ocsp >/dev/null
fi
